# matchbox/assets/coreos/oem-cloudinit.service
[Service]
Type=oneshot
ExecStart=/usr/bin/coreos-cloudinit --oem=ec2-compat

[Install]
WantedBy=multi-user.target